# ip.raspiotg/templates/systemd.service.j2
[Unit]
Description={{ item.description }}

[Service]
Type=simple
ExecStart={{ item.exec }}
Restart=always

[Install]
WantedBy=multi-user.target
