# ip.weewx/templates/systemd.service.j2
[Unit]
Description=Weather station hardware utility

[Service]
Type=simple
ExecStart={{ weewx_root }}/bin/weewxd
Restart=always

[Install]
WantedBy=multi-user.target
